# hxpctf2020/kernel-rop-bf9c106d45917343/kernel-rop/initramfs/etc/inittab
::sysinit:/etc/init.d/rcS
::once:-sh -c 'setuidgid 1000 sh; poweroff'
